# raspberrypi/huffle-shuffle.service
[Unit]
Description=HuffleShuffle Raspberry Pi Daemon
After=network-online.target
Wants=network-online.target
StartLimitInterval=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/huffle-shuffle/raspberrypi
ExecStart=/usr/bin/env bash /home/pi/huffle-shuffle/raspberrypi/startup.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huffle-shuffle

# Environment variables (will be loaded from .env file)
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=no
ReadWritePaths=/home/pi/huffle-shuffle
ReadWritePaths=/home/pi/.huffle

[Install]
WantedBy=multi-user.target
